# bingwallpaper/bingwallpaper-run.timer
[Unit]
Description=auto run BingWallpaper script when login

[Timer]
OnBootSec=60s
RemainAfterElapse=false

[Install]
WantedBy=timers.target